# deploy/systemd/fluxpoint.service
[Unit]
Description=Fluxpoint Next.js application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fluxpoint
Group=fluxpoint
WorkingDirectory=/var/www/fluxpoint
Environment=NODE_ENV=production
Environment=PORT=3021
Environment=HOST=127.0.0.1
Environment=HOSTNAME=127.0.0.1
EnvironmentFile=/var/www/fluxpoint/.env.production
ExecStart=/usr/bin/node /var/www/fluxpoint/.next/standalone/server.js
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
SyslogIdentifier=fluxpoint

[Install]
WantedBy=multi-user.target
